# supervisor.conf/apirest.ini
[program:apiRest]
command=npm start
autostart=true
autorestart=true
process_name=apiRest
startsecs=5
stopsignal=INT
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
user=root
directory=/api/transactions/ 
